# macos/Library/LaunchAgents/com.github.iakunin.dotfiles.gpgconf.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!-- Setting GPG_TTY and SSH_AUTH_SOCK env-variables globally (for all MacOS apps) -->
<!-- For more info see: https://stackoverflow.com/a/26586170/3456163 -->

<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.github.iakunin.dotfiles.gpgconf</string>

        <key>ProgramArguments</key>
        <array>
            <string>sh</string>
            <string>-c</string>
            <string>
                launchctl setenv GPG_TTY "$(tty)"
                launchctl setenv SSH_AUTH_SOCK $(/opt/homebrew/bin/gpgconf --list-dirs agent-ssh-socket)
                /opt/homebrew/bin/gpgconf --launch gpg-agent
            </string>

        </array>

        <key>RunAtLoad</key>
        <true/>

<!--        <key>StandardErrorPath</key>-->
<!--        <string>/Users/iakunin/bin/log/gpgconf.err</string>-->

<!--        <key>StandardOutPath</key>-->
<!--        <string>/Users/iakunin/bin/log/gpgconf.out</string>-->

    </dict>
</plist>
